# init/meetme.service
[Unit]
Description=metme backend

[Service]
Type=simple
ExecStart=/usr/bin/meetme
Restart=always

[Install]
WantedBy=multi-user.target
